# roles/systweaks/files/systemd/workloads.slice
[Unit]
Description=Slice for general workloads
Before=slices.target

[Slice]
AllowedCPUs=0-2
AllowedMemoryNodes=0
MemoryAccounting=yes
CPUAccounting=yes
IOAccounting=yes
Delegate=yes
