# roles/motra/files/systemd/isolated.slice
[Unit]
Description=Slice for isolated plc-server single device
Before=slices.target

[Slice]
AllowedCPUs=3
AllowedMemoryNodes=0
